# debian/shared-hosting.cron.d
0 */12 * * * root letsencrypt renew --quiet --non-interactive
0 */12 * * * root systemctl reload nginx
