# debian/shared-hosting.cron.d
0 */12 * * * root systemctl reload nginx
